# pleo-antaeus-scheduler/crontab
0 0 1 * * /scripts/trigger_pending.sh
*/5 * * * * /scripts/trigger_retry.sh
