# pleo-antaeus-scheduler/crontab
0 0 1 * * /scripts/triggerProcessing.sh
