# pleo-antaeus-scheduler/crontab
0 0 1 * * /scripts/triggerPending.sh
*/5 * * * * /scripts/triggerRetry.sh
